# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1366x768x24
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":1"
autorestart=true
priority=20

[program:chromium]
command=/usr/local/bin/start-browser.sh
environment=DISPLAY=":1"
autorestart=true
priority=30

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -forever -shared -rfbport 5901 -passwd %(ENV_VNC_PASSWORD)s
autorestart=true
priority=40

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 6901 localhost:5901
autorestart=true
priority=50

[program:cdp-forwarder]
command=/usr/bin/socat TCP-LISTEN:%(ENV_CDP_PORT)s,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:%(ENV_CHROME_DEBUG_PORT)s
autorestart=true
priority=60
